# deploy/systemd/harmony-resolver.service
[Unit]
Description=Harmony Resolver stack
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/harmony-resolver
EnvironmentFile=/etc/harmony-resolver/harmony-resolver.env
ExecStartPre=/usr/bin/docker compose -f compose.prod.yaml pull
ExecStart=/usr/bin/docker compose -f compose.prod.yaml up -d
ExecStop=/usr/bin/docker compose -f compose.prod.yaml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
